# ai-vault/launchagents/com.ai-empire.notes-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-empire.notes-ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/maurice/Documents/Dokumente – Mac mini von Maurice/New project/automation/scripts/run_notes_ingest.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/maurice/Documents/Dokumente – Mac mini von Maurice/New project/ai-vault/notes_ingest.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maurice/Documents/Dokumente – Mac mini von Maurice/New project/ai-vault/notes_ingest.err.log</string>
</dict>
</plist>
