# ai-vault/launchagents/com.ai-empire.telegram-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-empire.telegram-router</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/maurice/.openclaw/ai-vault/run-telegram-router.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/maurice/.openclaw/ai-vault/telegram_router.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maurice/.openclaw/ai-vault/telegram_router.err.log</string>
</dict>
</plist>
