# ai-vault/launchagents/com.ai-empire.shorts-revenue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-empire.shorts-revenue</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/maurice/Documents/New project/automation/scripts/run_shorts_revenue_once.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/maurice/Documents/New project/ai-vault/shorts_revenue.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maurice/Documents/New project/ai-vault/shorts_revenue.err.log</string>
</dict>
</plist>
