# ai-vault/launchagents/com.ai-empire.agentmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-empire.agentmonitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/maurice/Documents/New project/agent-dashboard/build/AgentMonitor.app/Contents/MacOS/AgentMonitor</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/maurice/.openclaw/ai-vault/agentmonitor.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maurice/.openclaw/ai-vault/agentmonitor.err.log</string>
</dict>
</plist>
