# scripts/com.aiempire.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiempire.telegram</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/maurice/AIEmpire-Core/telegram_bridge.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/maurice/AIEmpire-Core</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/maurice/AIEmpire-Core/logs/telegram_bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/maurice/AIEmpire-Core/logs/telegram_bridge_error.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
