# telegram/systemd/galaxia-orchestrator.service
[Unit]
Description=Galaxia OS Task Orchestrator - Background Worker
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/galaxia/orchestrator
EnvironmentFile=/opt/galaxia/orchestrator/.env
Environment="PYTHONUNBUFFERED=1"

ExecStart=/usr/bin/python3 /opt/galaxia/orchestrator/orchestrator.py

# Auto-restart on failure
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5

# Resource limits
MemoryLimit=2G
CPUQuota=100%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=galaxia-orchestrator

[Install]
WantedBy=multi-user.target
